# deploy/systemd/poly-market-arb-bot.service
[Unit]
Description=Polymarket Paper Arbitrage Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/poly-market-arb-bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/poly-market-arb-bot/.venv/bin/python scripts/paper_run.py
Restart=always
RestartSec=5
User=root
Group=root

[Install]
WantedBy=multi-user.target
